# roles/kibana/templates/kibana.upstart.j2
start on runlevel [2345]
respawn
setuid {{deploy_user}}
exec {{kibana_location}}/inst/bin/kibana